# AutoZilla/config/includes.chroot_after_packages/etc/profile.d/launcher.sh
sudo dhclient;
sleep 2.5
sudo bash /script/launch.sh
