# AutoZilla/config/includes.chroot_after_packages/etc/profile.d/launcher.sh
sudo cp /lib/live/mount/medium/preconfig.cfg .
sudo dmesg -D
sudo dhclient
sudo bash /script/launch.sh
