# AutoZilla/config/includes.chroot_after_packages/etc/profile.d/launcher.sh
sudo dhclient;
sudo bash /script/launch.sh
